# docker/cron/crontab
* * * * *  docker exec --user=application php php -d memory_limit=512M bin/console app:start-game
