# docker/cron/crontab
* * * * *  docker exec --user=www-data php php -d memory_limit=512M /var/www/symfony/bin/console app:start-game
